# usr/lib/systemd/system/systemd-udevd.service.d/etcmount.conf
[Unit]
# https://github.com/systemd/systemd/issues/13947
RequiresMountsFor=/etc
