# usr/lib/systemd/system/systemd-remount-fs.service.d/before-boot-writable.conf
[Unit]
# Make sure that at least one subvolume is mounted after the ro-remount
# of /, so that when local-fs.target is reached, the filesystem is writable
# again (boo#1156421). Otherwise, systemd-remount-fs.service might run after
# subvolume mounts, leaving / read-only (boo#1202000).
Before=boot-writable.mount
