# usr/lib/systemd/system/systemd-remount-fs.service.d/writableagain.conf
[Service]
# Make sure that at least one subvolume is mounted RW after the
# ro-remount of /, so that the filesystem is writable again as
# quickly as possible (boo#1156421).
ExecStart=/bin/sh -e -c 'if mountpoint -q /boot/writable; then mount -o remount,rw /boot/writable; fi'
