# usr/lib/systemd/system/systemd-journal-flush.service.d/afterlocalfs.conf
[Unit]
# https://bugzilla.opensuse.org/show_bug.cgi?id=1156421
After=local-fs.target
